# clearpath_robot/services/clearpath-manipulators.service
[Unit]
Description="Clearpath robot sub-service, launch all manipulators nodes"
PartOf=clearpath-robot.service
After=clearpath-robot.service

[Service]
User=administrator
Type=simple
ExecStart=/usr/sbin/clearpath-manipulators-start

[Install]
WantedBy=clearpath-robot.service
